# debian/postinst
#!/bin/sh
# postinst script for shaarli
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

case "$1" in
    configure)
        webservers="apache2, " #TODO: Add other servers
	    for webserver in $webservers; do
		    webserver=${webserver%,}
		    if [ "$webserver" = "apache2" ] && [ -f /etc/init.d/$webserver ] ; then
				# Enabling configuration is automatically added by dh_apache2
				# Fix owner and group for data, cache, pagecache and tmp folders
				chown www-data:www-data /var/lib/shaarli/data /var/cache/shaarli/cache /var/cache/shaarli/pagecache /var/cache/shaarli/tmp
		    #else
			#    # TODO: Handle other servers
		    fi
		    # Reload webserver in any case, configuration might have changed
		    # Redirection of 3 is needed because Debconf uses it and it might 
		    # be inherited by webserver. See bug #446324.
            if [ -f /etc/init.d/$webserver ] ; then
                if [ -x /usr/sbin/invoke-rc.d ]; then
                    invoke-rc.d $webserver reload 3>/dev/null || true
                else
                    /etc/init.d/$webserver reload 3>/dev/null || true
                fi
            fi
        done
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
